# backend/start_katalyzer
#!/bin/bash
#
#  katalyzer,v 0.10 2008/05/04 kovacik
#
# katalyzer:        Starts the katalyzer measurement 
#
# chkconfig: - 41 89
# description: This is a daemon which discovers routers on the local subnet.
# starting program /usr/local/bin/katalyzer
# processname: katalyzer
# config: /etc/katalyzer/my_config.conf

# Source function library.
. /etc/rc.d/init.d/functions

[ -f /etc/katalyzer/my_config.conf ] && . /etc/katalyzer/my_config.conf

start() {
	if [ $UID -ne 0 ]; then
                RETVAL=1
#                echo "only root can do that"
#		exit
		failure
        else
	        echo "Starting katalyzer: "
	        result=`ps -ef|grep /katalyzer |grep -v grep `

		if [[ $result == "" ]]; then /usr/local/bin/katalyzer -c /etc/katalyzer/my_config.conf > /dev/null &
		echo "[ OK ]";
		else 	
			echo "process: \"$result\" is running."
			echo "cannot start two such processes."
		fi;
		RETVAL=$?
                [ $RETVAL -eq 0 ] && touch /var/lock/subsys/katalyzer
	fi;
	echo
	return $RETVAL
}

stop() {
	if [ $UID -ne 0 ]; then
                RETVAL=1
                echo "only root can do that"
		exit
        else
		echo "Shutting down katalyzer. "
#	        kill `ps -ef|grep /katalyzer |grep -v grep |awk '{ print $2; }';`	
		killproc katalyzer
                [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/katalyzer
		RETVAL=$?
        fi;
	echo
	return $RETVAL
}

restart() {
	stop
	start
}

RETVAL=0

# See how we were called.
case "$1" in
  start)
	start
	;;
  stop)
	stop
	;;
  status)
	status katalyzer
	;;
  restart|reload)
	restart
	;;
  *)
	echo $"Usage: $0 {start|stop|status|restart}"
	exit 1
esac

exit $?
